# systemd/.config/systemd/user/gemini.service
[Unit]
Description=Gemini CLI Service
After=network.target

[Service]
ExecStart=/usr/local/bin/gemini

[Install]
WantedBy=default.target
