# systemd/.config/systemd/user/update-cover.loop.service
[Unit]
Description=Arrancar update-cover.loop.service Imagenen/cover de VLC.org.flatpak, etc

[Service]
ExecStart=/home/diego/dotfiles/home/update-cover-loop.sh
Type=oneshot

[Install]
WantedBy=default.target
